# etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

cd /etc/nginx
exec nginx
